# base/bin/rc.local
#!/bin/sh -e

#Fan
[ -f /sys/class/fan/mode ] && chmod 777 /sys/class/fan/mode
[ -f /sys/class/fan/level ] && chmod 777 /sys/class/fan/level
[ -f /sys/class/fan/enable ] && chmod 777 /sys/class/fan/enable
